# supervisor-jobmanager.ini
[supervisord]

[program:flink-jobmanager]
command = start-flink.sh jobmanager
autostart = true
autorestart = true
