# supervisor-taskmanager.ini
[supervisord]

[program:flink-taskmanager]
command = start-flink.sh taskmanager
autostart = true
autorestart = true
