# supervisor-jobmanager.ini
[supervisord]

[program:flink-jobmanager]
command = start-flink.sh jobmanager
autostart = true
autorestart = false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[eventlistener:flink-jobmanager_exit]
command=/usr/bin/kill-supervisor.py
process_name=flink-jobmanager
events=PROCESS_STATE_FATAL,PROCESS_STATE_EXITED,PROCESS_STATE_STOPPED

